# Software/Paneldactyl/supervisord.ini
[supervisord]
nodaemon=true


[program:panel]
user=nginx
group_name=nginx
autostart=true
autorestart=true
command=php /home/container/artisan queue:work --queue=high,standard,low --sleep=3 --tries=3
startsecs=1
startretries=3

[program:cron]
command=/usr/sbin/crond -f
user=root
autostart=true
autorestart=true